# ops/systemd/strikecore-backend.service
[Unit]
Description=Strikecore FastAPI backend (Hermes-style dashboard)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/argus-intelligence/strikecore
EnvironmentFile=%h/argus-intelligence/strikecore/.env
ExecStart=%h/argus-intelligence/strikecore/strikecore/bin/python3 -m uvicorn web.backend.app:app \
    --host 127.0.0.1 --port 8765 --proxy-headers --log-level info
Restart=on-failure
RestartSec=3
StandardOutput=append:%h/.strikecore/logs/backend.log
StandardError=append:%h/.strikecore/logs/backend.log

[Install]
WantedBy=default.target
